# setup/mimick.service
[Unit]
Description=Mimick Background Service
After=network.target network-online.target graphical-session.target
Wants=network-online.target

[Service]
# Runs as the current user
Type=simple
Environment=DISPLAY=:0
# Adjust the ExecStart to point to your python environment and src/main.py
# Example: ExecStart=/usr/bin/python3 /home/user/app/src/main.py --no-tray
ExecStart=%h/.local/bin/mimick --no-tray
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
